# services/auth/b25-auth.service
[Unit]
Description=B25 Authentication Service
Documentation=https://github.com/b25/auth-service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=b25
Group=b25
WorkingDirectory=/opt/b25/auth
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/b25/auth/src/server.js

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=b25-auth

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/b25
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
